# packaging/systemd/claude-contextd.service.in
[Unit]
Description=Claude Context daemon
After=network.target

[Service]
ExecStart=@@BIN_PATH@@
Restart=on-failure
RestartSec=2
Environment=HOME=@@HOME@@
StandardOutput=append:@@LOG_PATH@@
StandardError=append:@@LOG_PATH@@

[Install]
WantedBy=default.target
