# packaging/systemd/lm-semantic-search-daemon.service.in
[Unit]
Description=Claude Context daemon
Documentation=https://github.com/agoodkind/lm-semantic-search
After=network.target

[Service]
ExecStart=@@BIN_PATH@@
Restart=on-failure
RestartSec=2
Environment=HOME=@@HOME@@

[Install]
WantedBy=default.target
